# config/supervisord/supervisord.conf
# config/supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:meowcoin]
command=/usr/bin/meowcoind -conf=/home/meowcoin/.meowcoin/meowcoin.conf
user=meowcoin
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10
startretries=3
stopwaitsecs=60
stopasgroup=true
killasgroup=true

[program:fail2ban]
command=/bin/bash -c "if [ '${ENABLE_FAIL2BAN}' = 'true' ]; then /usr/bin/fail2ban-server -f; else sleep infinity; fi"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
startretries=3

[program:prometheus-exporter]
command=/bin/bash -c "if [ '${ENABLE_METRICS}' = 'true' ]; then /usr/local/bin/start-exporter.sh; else sleep infinity; fi"
user=meowcoin
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=5
startretries=3

[program:crond]
command=/bin/bash -c "if [ '${ENABLE_BACKUPS}' = 'true' ]; then crond -f -L /var/log/cron.log; else sleep infinity; fi"
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0